# mm-calendar.timer
[Unit]
Description=Run {{ SERVICE_NAME }} daily

[Timer]
OnCalendar=Mon..Fri 10:35
Persistent=true

[Install]
WantedBy=timers.target
